# daemon/oc-bridge.service
# systemd --user unit para manter o daemon sempre vivo.
# Instalar:
#   mkdir -p ~/.config/systemd/user
#   cp daemon/oc-bridge.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now oc-bridge
#   loginctl enable-linger pedrochagas   # roda mesmo sem sessão aberta
# Logs: journalctl --user -u oc-bridge -f

[Unit]
Description=oc-bridge (AITerminalControl) — ponte OpenClaude <-> Supabase
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/pedrochagas/Documentos/Projetos/AITerminalControl/daemon
# Serviço --user tem PATH mínimo: incluímos node (nvm), pnpm e bun explicitamente.
Environment=PATH=/home/pedrochagas/.nvm/versions/node/v24.14.1/bin:/home/pedrochagas/.local/share/pnpm:/home/pedrochagas/.bun/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/pedrochagas/.local/share/pnpm/pnpm start
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
